# scmclient/package/darwin/postinstall
#!/bin/bash

# --- 1. Variables ---
PLIST_LABEL="io.openscm.scmclient"
PLIST_PATH="/Library/LaunchDaemons/${PLIST_LABEL}.plist"
BINARY_PATH="/usr/local/bin/scmclient"

# --- 2. Fix Permissions (Mandatory for LaunchDaemons) ---
# If these are wrong, launchd will refuse to load the service
/usr/sbin/chown root:wheel "$PLIST_PATH"
/bin/chmod 644 "$PLIST_PATH"
/usr/sbin/chown root:wheel "$BINARY_PATH"
/bin/chmod 755 "$BINARY_PATH"

# --- 3. Restart the Service ---
echo "Stopping existing $PLIST_LABEL service (if running)..."

# Unload stops the current process and removes it from the service list
# We use '|| true' so the script doesn't crash on a fresh install
/bin/launchctl unload "$PLIST_PATH" 2>/dev/null || true

# Give the OS a second to release the file handles
/bin/sleep 1

echo "Starting new $PLIST_LABEL service..."
# Load -w tells macOS to start it now and keep it enabled across reboots
/bin/launchctl load -w "$PLIST_PATH"

exit 0


