# scmclient/package/darwin/postinstall
#!/bin/bash
# Post-install script — runs AFTER the Payload (binary + plist) is extracted.
# The service was already stopped by preinstall, so we only need to fix
# permissions and start the new binary.

LABEL="io.openscm.scmclient"
PLIST_PATH="/Library/LaunchDaemons/${LABEL}.plist"
BINARY_PATH="/usr/local/bin/scmclient"
LOG_DIR="/var/log/openscm"

# --- 1. Ensure directories exist ---
/bin/mkdir -p "$LOG_DIR"
/bin/chmod 755 "$LOG_DIR"

# --- 2. Fix permissions (LaunchDaemons WILL refuse to load with loose perms) ---
/usr/sbin/chown root:wheel "$BINARY_PATH"
/bin/chmod 755 "$BINARY_PATH"

/usr/sbin/chown root:wheel "$PLIST_PATH"
/bin/chmod 644 "$PLIST_PATH"

# --- 3. Enable + start the service ---
echo "Starting $LABEL service..."

# Enable ensures the service starts even if it was previously disabled
# (e.g. via 'launchctl disable' or 'launchctl unload -w' in an older install).
/bin/launchctl enable system/"$LABEL" 2>/dev/null || true

# Bootstrap loads the plist into the system domain and starts the process.
if ! /bin/launchctl bootstrap system "$PLIST_PATH"; then
    echo "Warning: launchctl bootstrap failed, attempting legacy load..."
    /bin/launchctl load -w "$PLIST_PATH" 2>/dev/null || true
fi

exit 0
