# scmclient/package/darwin/postinstall
#!/bin/bash
# Post-install script — runs AFTER the Payload (binary + plist) is extracted.
# We do a defensive bootout here before bootstrap: with KeepAlive:true, launchd
# may have re-bootstrapped the service between preinstall and now, which would
# cause bootstrap to fail and leave the old binary running.

LABEL="io.openscm.scmclient"
PLIST_PATH="/Library/LaunchDaemons/${LABEL}.plist"
BINARY_PATH="/usr/local/bin/scmclient"
LOG_DIR="/var/log/openscm"

# --- 1. Ensure directories exist ---
/bin/mkdir -p "$LOG_DIR"
/bin/chmod 755 "$LOG_DIR"

# --- 2. Fix permissions (LaunchDaemons WILL refuse to load with loose perms) ---
/usr/sbin/chown root:wheel "$BINARY_PATH"
/bin/chmod 755 "$BINARY_PATH"

/usr/sbin/chown root:wheel "$PLIST_PATH"
/bin/chmod 644 "$PLIST_PATH"

# --- 3. Stop any running instance before loading the new binary ---
# KeepAlive:true means launchd may have re-bootstrapped the service between
# preinstall and now. If we call bootstrap while it is already in the domain
# the call fails and the old binary keeps running.  Bootout here is
# unconditional and harmless if the service is not currently loaded.
echo "Stopping $LABEL (if running)..."
/bin/launchctl bootout system/"$LABEL" 2>/dev/null || true
/bin/sleep 2

# --- 4. Enable + start the service ---
echo "Starting $LABEL service..."

# Enable ensures the service starts even if it was previously disabled
# (e.g. via 'launchctl disable' or 'launchctl unload -w' in an older install).
/bin/launchctl enable system/"$LABEL" 2>/dev/null || true

# Bootstrap loads the plist into the system domain and starts the process.
if ! /bin/launchctl bootstrap system "$PLIST_PATH"; then
    echo "Warning: launchctl bootstrap failed, attempting legacy load..."
    /bin/launchctl load -w "$PLIST_PATH" 2>/dev/null || true
fi

exit 0
