# scmserver/debian/postinst
#!/bin/sh
set -e

APP_NAME="scmserver"
LOG_DIR="/var/log/openscm"

if [ "$1" = "configure" ]; then
    # 1. Create Log Directory
    if [ ! -d "$LOG_DIR" ]; then
        mkdir -p "$LOG_DIR"
    fi
    chmod 755 "$LOG_DIR"

    # 2. Reload and Start Service
    if [ -d /run/systemd/system ]; then
        systemctl --system daemon-reload >/dev/null || true
        systemctl enable "$APP_NAME.service" >/dev/null || true
        systemctl restart "$APP_NAME.service" >/dev/null || true
    fi

    # 3. Success Message
    echo ""
    echo "============================================================"
    echo "       $APP_NAME INSTALLATION COMPLETE"
    echo "============================================================"
    echo ""
fi

exit 0
